# archive/ubuntu/modules/base/templates/profile.d/prompt.sh.erb
#######################
## Managed by Puppet ##
#######################
export PS1='[\[\033[41;1m\]<%= @cluster %>\[\033[0m\]] \u@\h:\w\$ '
