# archive/ubuntu/modules/base/templates/profile.d/titlebar.sh.erb
#!/bin/bash
#(@) Joe DeCello
#######################
## Managed by Puppet ##
#######################
TITLE="[<%= @pop %>] $USER@<%= @fqdn %> [<%= @ipaddress %>]"
PROMPT_COMMAND='echo -ne "\033]0; $TITLE \007"'
