# archive/ubuntu/modules/base2011/files/profile.d/prompt.sh_test
export PS1='[\[\033[43;30m\]TEST\[\033[0m\]] \u@\h:\w\$ '
