# archive/ubuntu/modules/base2011/files/profile.d/prompt.sh_preview
export PS1='[\[\033[41;1m\]PREVIEW\[\033[0m\]] \u@\h:\w\$ '
