# archive/ubuntu/modules/base2011/files/profile.d/prompt.sh_qa
export PS1='[\[\033[43;30m\]QA\[\033[0m\]] \u@\h:\w\$ '
