# archive/ubuntu/modules/base2011/files/profile.d/prompt.sh_dev
export PS1='[\[\033[0;32m\]DEV\[\033[0m\]] \u@\h:\w\$ '
