# archive/ubuntu/modules/base2011/files/profile.d/prompt.sh_kx
export PS1='[\[\033[43;30m\]ex\[\033[0m\]] \u@\h:\w\$ '
