# archive/ubuntu/modules/base2011/files/profile.d/prompt.sh_DEFAULT
export PS1='[\[\033[41;1m\]WTF\[\033[0m\]] \u@\h:\w\$ '
