# archive/centos/modules/base/files/profile.d/prompt.sh_stage
export PS1='[\[\033[43;30m\]STAGE\[\033[0m\]] \u@\h:\w\$ '
