# archive/centos/modules/base/files/profile.d/prompt.sh_prod
export PS1='[\[\033[41;1m\]PROD\[\033[0m\]] \u@\h:\w\$ '
