# archive/centos/modules/base/files/profile.d/prompt.sh_DEFAULT
export PS1='[\[\033[41;1m\]example[\033[0m\]] \u@\h:\w\$ '
