# archive/centos/modules/base/templates/profile.d/titlebar.sh.erb
#!/bin/bash
TITLE="[<%= @environment %>] $USER@<%= @fqdn %> [<%= @ipaddress %>]"
PROMPT_COMMAND='echo -ne "\033]0; $TITLE \007"'
